# os_background_service/com.homeassist.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service identifier -->
    <key>Label</key>
    <string>com.homeassist.watchdog</string>

    <!-- Script to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/morgannstuart/Desktop/HomeAssistV3/os_background_service/watchdog.sh</string>
    </array>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/morgannstuart/Desktop/HomeAssistV3</string>

    <!-- Start on login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Always keep running -->
    <key>KeepAlive</key>
    <true/>

    <!-- Faster restart on failure -->
    <key>ThrottleInterval</key>
    <integer>5</integer>

    <!-- Log output -->
    <key>StandardOutPath</key>
    <string>/Users/morgannstuart/Desktop/HomeAssistV3/logs/watchdog_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/morgannstuart/Desktop/HomeAssistV3/logs/watchdog_stderr.log</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/morgannstuart</string>
    </dict>

    <!-- Lower priority than main assistant -->
    <key>Nice</key>
    <integer>5</integer>

    <!-- Background process type -->
    <key>ProcessType</key>
    <string>Background</string>

    <!-- Low priority I/O is fine for watchdog -->
    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
