# os_background_service/com.homeassistv2.assistant.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service identifier -->
    <key>Label</key>
    <string>com.homeassistv2.assistant</string>

    <!-- Script to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/morgannstuart/Desktop/HomeAssistV2/os_background_service/run_assistant.sh</string>
    </array>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/morgannstuart/Desktop/HomeAssistV2</string>

    <!-- Start on login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Unconditional keep alive - always restart regardless of exit code -->
    <key>KeepAlive</key>
    <true/>

    <!-- Reduced throttle interval for faster recovery (5 seconds) -->
    <key>ThrottleInterval</key>
    <integer>5</integer>

    <!-- Fallback watchdog: also trigger every 60 seconds as a safety net -->
    <!-- This ensures the service starts even if KeepAlive somehow fails -->
    <key>StartInterval</key>
    <integer>60</integer>

    <!-- Log output -->
    <key>StandardOutPath</key>
    <string>/Users/morgannstuart/Desktop/HomeAssistV2/logs/assistant_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/morgannstuart/Desktop/HomeAssistV2/logs/assistant_stderr.log</string>

    <!-- Environment variables (basic ones - rest come from .env) -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/morgannstuart</string>
    </dict>

    <!-- Nice priority (0 = normal, negative = higher priority) -->
    <key>Nice</key>
    <integer>-5</integer>

    <!-- Process type for macOS resource management -->
    <!-- Interactive ensures the process gets CPU/memory priority -->
    <key>ProcessType</key>
    <string>Interactive</string>

    <!-- Ensure I/O operations are not deprioritized -->
    <key>LowPriorityIO</key>
    <false/>

    <!-- Prevent launchd from killing on memory pressure -->
    <key>AbandonProcessGroup</key>
    <false/>
</dict>
</plist>
